# MorgAI/shared-host-infra-scaffold/infra/systemd/mortgage-backend-idle.service
[Unit]
Description=Idle-stop watcher run for mortgage-backend

[Service]
Type=oneshot
ExecStart=/opt/projects/mortgage-knowledge-assistant/infra/scripts/idle_stop_watcher.sh
