# MorgAI/shared-host-infra-scaffold/infra/systemd/mortgage-backend.socket
[Unit]
Description=Socket for on-demand activation of mortgage-backend

[Socket]
ListenStream=127.0.0.1:8001
Accept=no

[Install]
WantedBy=sockets.target
